# ramfs/overlays/BUSYBOX/etc/inittab
::sysinit:/bin/busybox chmod -R +x /etc/init.d/
::sysinit:/etc/init.d/rcS

::shutdown:/bin/umount -a -r
::shutdown:/bin/sync